# etc/systemd/system/containerd.service.d/override.conf
[Unit]
Requires=var-lib-containerd.mount

[Service]
Restart=always
RestartSec=5s
