# etc/systemd/system/mnt-storage-mountpoints.service.d/containerd.conf
[Service]
ExecStart=/usr/bin/env mkdir -p /mnt/storage/var/lib/containerd
