# etc/systemd/system/mnt-storage.automount
[Unit]
Description=Automount external storage

[Automount]
Where=/mnt/storage

[Install]
WantedBy=multi-user.target
